# deploy/services/falcon-orchestrator.service
[Unit]
Description=Falcon Trade Orchestrator
After=network-online.target falcon-dashboard.service
Wants=network-online.target
Requires=falcon-dashboard.service

[Service]
Type=simple
User=falcon
Group=falcon
WorkingDirectory=/home/ospartners/src/falcon

# Environment
EnvironmentFile=/opt/falcon/.env
Environment="PYTHONUNBUFFERED=1"
Environment="FALCON_ENV=production"

# Start command
ExecStart=/home/ospartners/src/falcon/backtest/bin/python3 -u /home/ospartners/src/falcon/run_orchestrator.py

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=falcon-orchestrator

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ospartners/src/falcon
ReadWritePaths=/var/lib/falcon

# Resource limits
MemoryMax=2G
CPUQuota=75%

[Install]
WantedBy=multi-user.target
